# device-unbound-el10/src/main/resources/usr/lib/systemd/system/unbound-postdiscovery.service
# systemd service for managing all unbound postconfiguration.

[Unit]
Description=Unbound Configuration Postdiscovery
After=unbound.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sequence -s daemon.info -b /usr/libexec device/autodiscovery/unbound.d start
ExecReload=/usr/bin/sequence -s daemon.info -b /usr/libexec device/autodiscovery/unbound.d reload
RemainAfterExit=on

[Install]
WantedBy=multi-user.target

